# desktop/resources/speech2prompt.service
[Unit]
Description=Speech2Prompt Desktop - Voice-to-keyboard bridge
Documentation=https://github.com/yourusername/speech2prompt
After=graphical-session.target bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
ExecStart=/usr/local/bin/speech2prompt-desktop
Restart=on-failure
RestartSec=5
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/%U

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.config/speech2prompt %h/.local/share/speech2prompt
PrivateTmp=yes

[Install]
WantedBy=default.target
